# examples/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.streaming-catalog.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/streaming-catalog</string>
        <string>sync</string>
        <string>--collect-first</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/streaming-catalog</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/streaming-catalog.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/streaming-catalog.err</string>
</dict>
</plist>
